# systemd/minhos-sierra-client.service
[Unit]
Description=MinhOS Sierra Chart Client Service
Documentation=https://github.com/minh_os_v3
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=colindo
Group=colindo
WorkingDirectory=/home/colindo/Sync/minh_v3
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/colindo/Sync/minh_v3"
Environment="BRIDGE_HOSTNAME=cthinkpad"
Environment="BRIDGE_PORT=8765"

# Use the main sierra_client.py with integrated optimizations
ExecStart=/usr/bin/python3 /home/colindo/Sync/minh_v3/minhos/services/sierra_client.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitBurst=5

# Resource limits
LimitNOFILE=65535
LimitCORE=infinity

# Logging
StandardOutput=journal
StandardError=journal

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target